# dist/upstart.example.conf
#/etc/init/echodemo.conf
#log file must exist if its path on config is not writeable "touch /var/log/echodemo.log;chown echodemo:echodemo /var/log/echodemo.log"
description "ECHO Demo webapp upstart daemon"
author "Gerónimo Oñativia <geronimox@gmail.com>"

start on runlevel [2345]
stop on runlevel [!2345]

setuid echodemo
setgid echodemo
chdir /home/echodemo/app

console none
respawn

exec /home/echodemo/app/echodemo -config /home/echodemo/app/echodemo.ini